# sandboxes/research-agent/netns/agent-allowlist.timer
[Unit]
Description=Refresh SANAD research-agent egress allowlist every 60s

[Timer]
OnBootSec=30s
OnUnitActiveSec=60s
AccuracySec=5s
Unit=agent-allowlist.service

[Install]
WantedBy=timers.target
